# systemd/idolhub-mcp.service.template
[Unit]
Description=idolhub MCP stdio Server
Documentation=https://github.com/primit1v0/idolhub-bot
After=network.target

[Service]
Type=simple
User=@IDOLHUB_USER@
Group=@IDOLHUB_GROUP@
WorkingDirectory=@IDOLHUB_DIR@
EnvironmentFile=@IDOLHUB_ENV_FILE@
ExecStart=@IDOLHUB_DIR@/.venv/bin/python main.py mcp
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=idolhub-mcp

[Install]
WantedBy=default.target
